# server-setup/conf/supervisor.conf.tpl
[program:{program_name}]
directory=/opt/{project_folder}/dashboard-metrics/
command=/opt/{project_folder}/bin/start.sh
user=developer
autostart=true
redirect_stderr=true
